# lib/itamae/plugin/recipe/consul/service/templates/etc/systemd/system/consul.service.erb
[Unit]
Description=consul agent
Requires=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-<%= @environment_file %>
Environment=GOMAXPROCS=2
Restart=on-failure
ExecStart=<%= @bin_path %> agent $OPTIONS -config-dir=/etc/consul.d -data-dir=<%= @data_dir %>
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
